# rainbow.service
[Unit]
Description=Example systemd service.

[Service]
Type=exec
ExecStart=/usr/local/bin/rainbow

[Install]
WantedBy=multi-user.target
